# linux/embedded_rootfs/etc-files/inittab
::sysinit:/sbin/rc boot
::respawn:-/bin/cav_sh /bin/sh
